# devenv-jessie/rc.local

# Start services
# TODO: find a better way, System V ?

service mesos-slave start
service marathon start

exit 0
