# jessie-vagrant/slave-files/rc.local
#!/usr/bin/env bash

# Start services
# TODO: find a better way, System V ?

service mesos-slave start
service haproxy start

exit 0
