# scripts/master-files/rc.local
#!/usr/bin/env bash

# Start services
# TODO: find a better way, System V ?

service mesos-master start
service marathon start
service chronos start
service haproxy start

exit 0
